# docker/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:symfony]
command=symfony server:start
autostart=true
autorestart=true